# fedora-prime.service
[Unit]
Description=Enable intel GPU on boot
Before=gdm.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/fedora-prime-select intel

